# sync-ipv6-dns.timer
[Unit]
Description=Run sync-ipv6-dns daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d
Unit=sync-ipv6-dns.service

[Install]
WantedBy=timers.target
